# crontab
* * * * * cd /var/www/html $$ php artisan reminders:send >> /dev/null 2>&1
